# files/etc/crontabs/mycron
# CRON TABLE
0 0 * * * root /usr/local/bin/geoip_database.sh
